# stage8/10-monit/files/update-motd.d/90-monit
#!/bin/sh
echo
echo "Services status: "
/usr/bin/monit report | sed -e 's/^/     /'
echo
